# deploy/gcp-vm-collector/cosyvoice-usage.service
[Unit]
Description=CosyVoice2 Korean SFT opt-in usage collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cosyvoice-usage
Group=cosyvoice-usage
WorkingDirectory=/opt/cosyvoice-usage
EnvironmentFile=/etc/cosyvoice-usage.env
ExecStart=/opt/cosyvoice-usage/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8092 --workers 1 --proxy-headers --no-access-log
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cosyvoice-usage

[Install]
WantedBy=multi-user.target
